# files-main/qubes-vpn-handler.timer
[Unit]
Description=Start vpn

[Timer]
OnBootSec=1

[Install]
WantedBy=timers.target
