# rw/config/rc.local
#!/bin/sh

# Need qvpn group for openvpn access through firewall
groupadd -r qvpn

# Add our openvpn service config and start
cp /rw/config/openvpn/openvpn-client.service /etc/systemd/system/
sync; systemctl daemon-reload
systemctl --no-block start openvpn-client.service
